# launch_agents/com.lky.arxive-agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.lky.arxive-agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/opt/anaconda3/bin/python3</string>
    <string>/Users/lky/Desktop/Agent/arxive_agent/scripts/run_scheduled.py</string>
    <string>--config</string>
    <string>/Users/lky/Desktop/Agent/arxive_agent/config/config.yaml</string>
  </array>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>8</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/lky/Desktop/Agent/arxive_agent/logs/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/lky/Desktop/Agent/arxive_agent/logs/launchd.err.log</string>
</dict>
</plist>
